# services/sh-rpi-venus/run
#!/bin/sh
exec 2>&1
exec /data/sh-rpi-venus/daemon.py